# v2/minecraft-stop.timer
[Unit]
Description=Stop Minecraft server

[Timer]
OnCalendar=*-*-* 04:45:00
Persistent=true

[Install]
WantedBy=timers.target
