# v2/minecraft-start.service
[Unit]
Description=Start Minecraft

[Service]
Type=oneshot
ExecStart=/bin/systemctl start minecraft
