# v2/minecraft-stop.service
[Unit]
Description=Stop Minecraft

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop minecraft