# v2/minecraft-backup.service
[Service]
Type=oneshot
ExecStart=/opt/minecraft/backup.sh