# v2/minecraft-backup.timer
[Unit]
Description=Backup Minecraft Server

[Timer]
OnCalendar=*-*-* 04:50:00
Persistent=true

[Install]
WantedBy=timers.target