# v2/minecraft-backup.service
[Unit]
Description=Backup Minecraft

[Service]
Type=oneshot
ExecStart=/opt/minecraft/backup.sh