# v2/minecraft-backup.timer
[Timer]
OnCalendar=*-*-* 04:50:00
Persistent=true